# server-config/systemd/transinfo-db-backup.service
[Unit]
Description=Transinfo DB backup

[Service]
Type=oneshot
User=postgres
Group=postgres
ExecStart=/usr/local/bin/transinfo-db-backup.sh
